# rootfs/service/rtorrent/run
#!/bin/bash

# Checking the lock
if [ -f /rtorrent/.session/rtorrent.lock ]; then
    rm -f /rtorrent/.session/rtorrent.lock
fi

# Turn to rTorrent
exec rtorrent
